# management _sysytem/web/rootfs/etc/services.d/seo-monitor-web/run
#!/usr/bin/with-contenv sh

cd /usr/src/app

# start our app
npm run start;